# files/etc/service/postgres/run
#!/bin/bash
exec setuser system ${PGENGINE}/postgres -D $PGDATA -c config_file=$PGCONF/postgresql.conf
